# k8s/clusters/production/infrastructure/democratic-csi/backup-target-saveconfig.timer
[Unit]
Description=Hourly backup of /etc/target/saveconfig.json

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
